# docs/telegram-bot.service
[Unit]
Description=Qwen3-TTS Telegram Bot Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=%u
Group=%g
WorkingDirectory=/path/to/tts-server
Environment="PYTHONPATH=/path/to/tts-server"
# Pick the interpreter for the family environment you want to run, for example:
#   /path/to/tts-server/.envs/qwen/bin/python
#   /path/to/tts-server/.envs/omnivoice/bin/python
ExecStart=/path/to/tts-server/.envs/qwen/bin/python -m telegram_bot
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=qwen3-tts-telegram

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/path/to/tts-server/.outputs /path/to/tts-server/.uploads /path/to/tts-server/.state

# Resource limits
LimitNOFILE=65536
MemoryMax=2G

# Health check (optional - uncomment if using systemd的健康检查)
# ExecStartPost=/bin/bash -c 'sleep 5 && systemctl is-active --quiet telegram-bot || exit 1'

[Install]
WantedBy=multi-user.target
